# templates/insightsupdate.cron.erb
#!/bin/bash

/usr/bin/insights-client >> /var/log/insights-status.log 2>&1
